# pkg/debian/prerm
#!/bin/sh -e

echo "Stopping processes"
systemctl stop skyd
systemctl disable skyd
echo "Stopped processes"
