# packages/debian/DEBIAN/postrm
#!/bin/sh

sudo userdel skytable
sudo rm -rf /var/opt/skytable

if [ -d /run/systemd/system ]; then
	systemctl --system daemon-reload >/dev/null || true
fi
if [ "$1" = "remove" ]; then
	if [ -x "/usr/bin/deb-systemd-helper" ]; then
		deb-systemd-helper mask skyd.service >/dev/null || true
	fi
fi

if [ "$1" = "purge" ]; then
	if [ -x "/usr/bin/deb-systemd-helper" ]; then
		deb-systemd-helper purge skyd.service >/dev/null || true
		deb-systemd-helper unmask skyd.service >/dev/null || true
	fi
fi
